# odyssey-x86j4125-code/files/etc/uci-defaults/99-adblock_config
# https://github.com/openwrt/packages/blob/master/net/adblock/files/README.md
uci -q batch << EOI
set adblock.global.adb_triggerdelay="60"
add_list adblock.global.adb_sources="adguard_tracking"
add_list adblock.global.adb_sources="android_tracking"
add_list adblock.global.adb_sources="andryou"
add_list adblock.global.adb_sources="anti_ad"
add_list adblock.global.adb_sources="antipopads"
add_list adblock.global.adb_sources="anudeep"
add_list adblock.global.adb_sources="bitcoin"
add_list adblock.global.adb_sources="cpbl"
add_list adblock.global.adb_sources="doh_blocklist"
add_list adblock.global.adb_sources="easylist"
add_list adblock.global.adb_sources="easyprivacy"
add_list adblock.global.adb_sources="firetv_tracking"
add_list adblock.global.adb_sources="games_tracking"
add_list adblock.global.adb_sources="hblock"
add_list adblock.global.adb_sources="lightswitch05"
add_list adblock.global.adb_sources="notracking"
add_list adblock.global.adb_sources="oisd_big"
add_list adblock.global.adb_sources="oisd_nsfw"
add_list adblock.global.adb_sources="oisd_small"
add_list adblock.global.adb_sources="openphish"
add_list adblock.global.adb_sources="phishing_army"
add_list adblock.global.adb_sources="reg_cn"
add_list adblock.global.adb_sources="reg_cz"
add_list adblock.global.adb_sources="reg_de"
add_list adblock.global.adb_sources="reg_es"
add_list adblock.global.adb_sources="reg_fi"
add_list adblock.global.adb_sources="reg_fr"
add_list adblock.global.adb_sources="reg_id"
add_list adblock.global.adb_sources="reg_it"
add_list adblock.global.adb_sources="reg_jp"
add_list adblock.global.adb_sources="reg_kr"
add_list adblock.global.adb_sources="reg_nl"
add_list adblock.global.adb_sources="reg_pl"
add_list adblock.global.adb_sources="reg_ro"
add_list adblock.global.adb_sources="reg_ru"
add_list adblock.global.adb_sources="reg_se"
add_list adblock.global.adb_sources="reg_vn"
add_list adblock.global.adb_sources="smarttv_tracking"
add_list adblock.global.adb_sources="spam404"
add_list adblock.global.adb_sources="stevenblack"
add_list adblock.global.adb_sources="stopforumspam"
add_list adblock.global.adb_sources="utcapitole"
add_list adblock.global.adb_sources="wally3k"
add_list adblock.global.adb_sources="whocares"
add_list adblock.global.adb_sources="winhelp"
add_list adblock.global.adb_sources="winspy"
commit adblock
EOI
