# odyssey-x86j4125-code/files/etc/uci-defaults/99-irqbalance_config
# https://openwrt.org/docs/guide-user/services/irqbalance
uci -q batch << EOI
set irqbalance.@irqbalance[0].enabled="1"
commit irqbalance
EOI
