# odyssey-x86j4125-code/files/etc/uci-defaults/99-uci-config
# https://openwrt.org/docs/guide-developer/uci-defaults#integrating_custom_settings
# https://openwrt.org/docs/guide-user/base-system/dhcp_configuration#replacing_dnsmasq_with_odhcpd_and_unbound
# https://github.com/openwrt/packages/blob/master/net/unbound/files/README.md
uci -q batch << EOI
set network.lan.ipaddr='192.168.2.1'
commit network
set unbound.@unbound[0].add_local_fqdn="3"
set unbound.@unbound[0].add_wan_fqdn="1"
set unbound.@unbound[0].dhcp_link="odhcpd"
set unbound.@unbound[0].dhcp4_slaac6="1"
set unbound.@unbound[0].unbound_control="1"
commit unbound
set dhcp.odhcpd.leasetrigger="/usr/lib/unbound/odhcpd.sh"
commit dhcp
set adblock.global.adb_triggerdelay="60"
commit adblock
set nlbwmon.@nlbwmon[0].netlink_buffer_size="1048576"
set nlbwmon.@nlbwmon[0].refresh_interval="10s"
commit nlbwmon
EOI
