# odyssey-x86j4125-code/files/etc/uci-defaults/99-unbound_config
# https://openwrt.org/docs/guide-user/base-system/dhcp_configuration#replacing_dnsmasq_with_odhcpd_and_unbound
# https://github.com/openwrt/packages/blob/master/net/unbound/files/README.md
uci -q batch << EOI
set unbound.@unbound[0].add_local_fqdn="3"
set unbound.@unbound[0].add_wan_fqdn="1"
set unbound.@unbound[0].dhcp_link="odhcpd"
set unbound.@unbound[0].dhcp4_slaac6="1"
set unbound.@unbound[0].unbound_control="1"
commit unbound
set dhcp.odhcpd.leasetrigger="/usr/lib/unbound/odhcpd.sh"
commit dhcp
EOI
