# odyssey-x86j4125-code/files/etc/uci-defaults/99_uci-config
# https://openwrt.org/docs/guide-developer/uci-defaults#integrating_custom_settings
uci -q batch << EOI
set network.lan.ipaddr='192.168.2.1'
commit network
EOI